# shadowplay.timer
[Timer]
OnCalendar=*:0/2
OnBootSec=1min

[Install]
WantedBy=timers.target
